# init.d/fortivpn.service
[Unit]
Description=Fortinet VPN Client Service

[Service]
ExecStart=/usr/bin/fortivpn
Restart=on-failure

[Install]
WantedBy=multi-user.target
